# docker/configs/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command = nginx -g "daemon off;"
redirect_stderr = true
autostart = true
